# ihp_ml_addon/rootfs/etc/services.d/ihp-ml/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e

# Read configuration from options.json (Home Assistant add-on config)
CONFIG_PATH=/data/options.json

if [ -f "$CONFIG_PATH" ]; then
    CONFIG_LOG_LEVEL=$(bashio::config 'log_level' 'info')
    CONFIG_MODEL_PATH=$(bashio::config 'model_persistence_path' '/data/models')
else
    CONFIG_LOG_LEVEL="info"
    CONFIG_MODEL_PATH="/data/models"
fi

# Set environment variables
export LOG_LEVEL="${CONFIG_LOG_LEVEL}"
export MODEL_PERSISTENCE_PATH="${CONFIG_MODEL_PATH}"

# Ensure model directory exists
mkdir -p "${MODEL_PERSISTENCE_PATH}"

bashio::log.info "Starting IHP ML Models Add-on..."
bashio::log.info "Log level: ${LOG_LEVEL}"
bashio::log.info "Model path: ${MODEL_PERSISTENCE_PATH}"

# Start the Flask application
cd /app
exec /opt/venv/bin/python -m infrastructure.api.server
